# automation/instance-factory/templates/systemd/agent-mail-http-bridge.service
[Unit]
Description=RemoteLab agent mail HTTP bridge
After=network-online.target remotelab.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory={{repoCheckoutPath}}
EnvironmentFile=/etc/remotelab/install-profile.env
ExecStart=/usr/bin/env node {{repoCheckoutPath}}/scripts/agent-mail-http-bridge.mjs
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
